# _o/shell_rc/sh_rc_scripts/profile.d/sh-prompt_2rows.sh
#/bin/sh

# sh ps1: 2 rows.
# /u2584=▄ /u2580=▀ /u258c=▌ /u2588=█
PS1='\[\033[0;32m\]▄[\u@\H \[\033[0;33m\]\w\[\033[0;32m\]]\n\[\033[0;32m\]▀\$\[\033[0m\] '
