# shell_rc/sh_rc_scripts/profile.d/_profile.d.sh
#/bin/sh

source ~/.profile.d/utility-alias.sh

source ~/.profile.d/sh-ps1.sh
